# roles/app/templates/zds-search-engine-index.timer.j2
[Unit]
Description=Search engine reindex new and updated contents

[Timer]
OnCalendar=*:30:00
Persistent=true

[Install]
WantedBy=timers.target
