# roles/app/templates/zds-es-index.service.j2
[Unit]
Description=Reindex ES Service

[Service]
Type=oneshot
User={{ appuser }}
Group={{ appuser }}
ExecStart={{ workdir }}/wrapper es_manager index_flagged
