# roles/app/templates/zds-watchdog.service.j2
[Unit]
Description=Zeste de Savoir - Watchdog
After=network.target

[Service]
User={{ appuser }}
Group={{ appuser }}
ExecStart={{ workdir }}/wrapper publication_watchdog
ExecStop=/bin/kill -s TERM $MAINPID
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
