# roles/app/templates/zds-es-index.timer.j2
[Unit]
Description=ES reindex flagged contents

[Timer]
OnCalendar=*:30:00
Persistent=true

[Install]
WantedBy=timers.target
