# roles/app/templates/zds.socket.j2
[Unit]
Description=ZdS Gunicorn socket

[Socket]
ListenStream={{ rundir }}/gunicorn.sock

[Install]
WantedBy=sockets.target
